# systemd-user/home-backup.timer
[Unit]
Description=Run home-backup once per day

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
