# systemd-user/random-wallpaper.timer
[Unit]
Description=Set random GNOME wallpaper once per hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
